# meta-jaspar/recipes-multimedia/gmrender-resurrect/gmrender-resurrect/gmrender-resurrect.service
[Unit]
Description=Headless UPnP Renderer
After=connman.service

[Service]
Type=forking
ExecStart=/usr/lib/gmediarender/gmediarender.sh

[Install]
WantedBy=multi-user.target
